# config/taskmonitor.service.in
###############################################################################
#  TaskMonitor version @GIT_SHA1@
###############################################################################
[Unit]
Description=TaskMonitor Service

[Service]
Type=simple
User=root
Group=root
ExecStart=@CMAKE_INSTALL_PREFIX@/sbin/taskmonitor -c /etc/taskmonitor.conf
WatchdogSec=10

[Install]
WantedBy=multi-user.target
